# w2c-letsencrypt
#!/bin/sh
#
# Copyright (c) Johannes Feichtner <johannes@web-wack.at>
# Released under the GNU GPLv3 License.
#
# chkconfig: on 99 99
# description: Letsencrypt for ESXi
#

# Paths for config and backup
BACKUP_CFG="/etc/renew.cfg.bak"
TARGET_CFG="/opt/w2c-letsencrypt/renew.cfg"

# If renew.cfg exists and is different from backup, update the backup
if [ -f "$TARGET_CFG" ]; then
    if [ ! -f "$BACKUP_CFG" ] || ! cmp -s "$TARGET_CFG" "$BACKUP_CFG"; then
        cp -f "$TARGET_CFG" "$BACKUP_CFG"
        echo "Backed up user renew.cfg to $BACKUP_CFG."
    fi
fi

# Restore user-modified renew.cfg from persistent backup if it exists
if [ -f "$BACKUP_CFG" ]; then
    cp -f "$BACKUP_CFG" "$TARGET_CFG"
    echo "Restored user renew.cfg from backup."
fi

export PATH=/sbin:/usr/sbin:/bin:/usr/bin

log() {
   echo "$@"
   logger -p daemon.info -t "$0" "$@"
}

for action in "$@"; do
   log "Running '${action}' action";

   case "$action" in
      start)
         /opt/w2c-letsencrypt/renew.sh
         ;;

      remove)
         sed -i '/\/opt\/w2c-letsencrypt/d' /var/spool/cron/crontabs/root
         sed -i '/acme-challenge/d' /etc/vmware/rhttpproxy/endpoints.conf
         /sbin/generate-certificates
         for s in /etc/init.d/*; do if $s | grep ssl_reset > /dev/null; then $s ssl_reset; fi; done
         ;;

      stop)
         true
         ;;

      status)
         true
         ;;

      install)
         true
         ;;

      upgrade)
         true
         ;;

      restart)
         "$0" stop
         "$0" start
         ;;

      *)
         echo "Usage: $(basename "$0") {start|stop|status|restart}"
         exit 1
   esac
done

exit 0
